# packages/inputmethods/fcitx/files/systemd/user/fcitx.service
[Unit]
Description=Fcitx Input Method
Requires=xorg.socket
After=xorg.socket

[Service]
ExecStart=/usr/bin/fcitx -D -r
