# packages/sys-devel/distcc/files/systemd/distccd.service
[Unit]
Description=distccd - distributed C/C++ compiler server

[Service]
User=nobody
EnvironmentFile=/etc/conf.d/distccd.conf
ExecStart=/usr/bin/distccd --listen ${LISTEN_ADDRESS} --port ${LISTEN_PORT} --allow ${ALLOW} --log-stderr --log-level ${LOG_LEVEL} --nice ${NICENESS} --no-detach --daemon ${EXTRA_OPTS}

[Install]
WantedBy=multi-user.target
