# packages/sys-devel/distcc/files/systemd/distccd.service
[Unit]
Description=distccd - distributed C/C++ compiler server

[Service]
User=nobody
ExecStart=/usr/bin/distccd --allow 0.0.0.0/0 --no-detach
