# packages/sys-devel/distcc/files/systemd/distccd.service
[Unit]
Description=distccd - distributed C/C++ compiler server

[Service]
User=nobody
EnvironmentFile=/etc/conf.d/distccd.conf
ExecStart=/usr/bin/distccd --listen ${ADDRESS} --port ${PORT} --allow ${ALLOW} --no-detach --daemon ${EXTRA_OPTS}
