# packages/x11-misc/synergy/files/systemd/user/synergy.service
[Unit]
Description=Synergy server

[Service]
ExecStart=/usr/bin/synergys --config ${HOME}/.config/synergy/config --debug NOTE --no-restart --no-daemon
