# packages/x11-misc/synergy/files/systemd/user/synergy@.service
[Unit]
Description=Synergy client @ host of %i

[Service]
ExecStart=/usr/bin/synergyc --no-daemon --restart --debug NOTE %i
