# packages/net-misc/libteam/files/systemd/teamd.service
[Unit]
Description=Team Daemon

[Service]
ExecStart=/usr/bin/teamd -U -D -f /etc/teamd.conf
