# packages/net-misc/wicd/files/systemd/wicd-suspend.service
[Unit]
Description=Wicd suspend hook
Before=sleep.target

[Service]
Type=simple
ExecStart=-/usr/share/wicd/daemon/suspend.py

[Install]
WantedBy=sleep.target

