# packages/net-misc/wicd/files/systemd/wicd-resume.service
[Unit]
Description=Wicd resume hook
After=suspend.target

[Service]
Type=simple
ExecStart=-/usr/share/wicd/daemon/autoconnect.py

[Install]
WantedBy=suspend.target
