# packages/media-sound/mpd/files/systemd/user/mpd.socket
[Socket]
#ListenStream=/run/mpd.socket
ListenStream=6600

[Install]
WantedBy=sockets.target
